# bin/service/ohttps.service
[Unit]
Description=Ohttps Webhook Service
After=network.target remote-fs.target nss-lookup.target

[Service]
User=root
Group=root
Type=forking
ExecStart=/usr/bin/ohttps -c /etc/ohttps.d/config.yaml -s

[Install]
WantedBy=multi-user.target