# bin/service/ohttps.service
[Unit]
Description=Ohttps Webhook Service
After=network.target

[Service]
User=root
Group=root
ExecStart=/usr/bin/ohttps -c /etc/ohttps.d/config.yaml
Restart=allways
RestartSec=5

[Install]
WantedBy=multi-user.target
